# systemd/weatherbot.service
[Unit]
Description=Weatherbot daemon responds to requests sent to cgi
After=network-online.target

[Service]
#Type=forking
ExecStart=/opt/weatherbotd
TimeoutStartSec=60
WorkingDirectory=/opt/
Restart=always
RestartSec=10
User=weatherbot
Group=weatherbot

[Install]
WantedBy=multi-user.target

